# ansible/roles/chicago-host/templates/prune-images.service.j2
[Unit]
Description=Prune dangling docker images
Wants=prune-images.timer

[Service]
ExecStart=/usr/bin/docker image prune --force --all
WorkingDirectory={{ chicon_home }}

[Install]
WantedBy=multi-user.target
