# ansible/roles/chicago-host/templates/rotate-creds.service.j2
[Unit]
Description=Runs credential rotator
Wants=rorate-creds.timer

[Service]
ExecStart=/opt/chicago/sbin/rotate-creds.sh
WorkingDirectory=/opt/chicago

[Install]
WantedBy=multi-user.target
